# scripts/assets/rc.local
#!/bin/bash
setcap 'cap_net_raw+p' /usr/bin/ping
getcap /usr/bin/ping
